# scripts/com.phd-os.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.phd-os.email-pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>
        <string>__SCRIPT_PATH__</string>
    </array>

    <!-- Run every 2 hours -->
    <key>StartInterval</key>
    <integer>7200</integer>

    <!-- Also run immediately when agent is loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Only run when user is logged in (LaunchAgent, not LaunchDaemon) -->

    <key>StandardOutPath</key>
    <string>/tmp/phd_os_pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/phd_os_pipeline.err</string>

    <!-- Prevent thrashing if the script crashes -->
    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
